# .config/autostart/goldendict.desktop
[Desktop Entry]
Name=goldendict
Exec=goldendict
Type=Application
X-GNOME-Autostart-enabled=true
EOS
